# screenlock-handler.service
[Unit]
Description='Lock screen'

[Service]
Restart=always
ExecStart=/usr/local/bin/screenlock-handler
